# systemd/eyesy-pd.service
[Unit] 
Description=Pd OSC hardware bridge for Eyesy

[Service] 
Type=simple 
WorkingDirectory=/home/we/Eyesy
ExecStart=/usr/bin/pd -alsamidi -nogui -noaudio -path /home/we/sidekick/pdexternals pd/main.pd
User=we
Group=we

[Install] 
WantedBy=multi-user.target
